# FW/nu4100/bootfix/Host/scripts/inuservice.conf
# inuservice

description	"inuservice"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn
respawn limit 10 5

exec /etc/init.d/inuservice start
